# groups/deploy/systemd/phoenix-listener.service
[Unit]
Description=Phoenix v2 WhatsApp listener (Baileys)
After=network-online.target phoenix-brain.service
Wants=phoenix-brain.service

[Service]
Type=simple
WorkingDirectory=%h/phoenix/wa-listener
EnvironmentFile=%h/phoenix/.env
EnvironmentFile=%h/phoenix/wa-listener/.env
ExecStart=/usr/bin/env npm run start
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/phoenix/logs/listener.log
StandardError=append:%h/phoenix/logs/listener.err

[Install]
WantedBy=default.target
